# django.cron
*/5 * * * * /usr/local/bin/python /usr/src/app/manage.py runcrons >> /usr/src/app/cron.log 2>&1
